# raspi/timekeeper.service
[Unit]
Description=TimeKeeper service

[Service]
ExecStart=/home/pi/TimeKeeper/service.py
WorkingDirectory=/home/pi/TimeKeeper
User=pi
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
